# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9001        ; (ip_address:port specifier, *:port for all iface)
username=pi                ; (default is no username (open server))
password=8888              ; (default is no password (open server))

[supervisord]
logfile=/home/tete/News-Weather/log/supervisord.log ; (main log file;default $CWD/supervisord.log)
childlogdir=/home/tete/News-Weather/log/subprocess_log
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:uwsgi]
command=/home/tete/News-Weather/venv/bin/uwsgi uwsgi.ini

[program:celery_worker]
command=/home/tete/News-Weather/venv/bin/celery -A tasks worker -l WARNING

[program:celery_beat]
command=/home/tete/News-Weather/venv/bin/celery -A tasks beat -l WARNING

;[include]
;files = relative/directory/*.ini
